# crontab
@reboot /home/pi/lanpartysign/startup.sh